# systemd/jack-connection-manager.service.in
[Unit]
Description=jack-connection-manager
After=network.target
[Install]
# There is no multi-user.target for a user service, therefore we use default.target
WantedBy=default.target

[Service]
Type=simple
ExecStart=@bin_dir@/jack-connection-manager -c /usr/local/etc/jack-connection-manager/connections.yml
# LimitRTPRIO=95
# LimitRTTIME=infinity
# LimitMEMLOCK=infinity