# systemd/jack-wire-puller.service.in
[Unit]
Description=jack-wire-puller
After=network.target
[Install]
# There is no multi-user.target for a user service, therefore we use default.target
WantedBy=default.target

[Service]
Type=simple
ExecStart=@bin_dir@/jack-wire-puller -c /usr/local/etc/jack-wire-puller/connections.yml
# LimitRTPRIO=95
# LimitRTTIME=infinity
# LimitMEMLOCK=infinity